# training/automation/systemd/codex-watcher.service
[Unit]
Description=Codex event watcher for D5/DeFi engine
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/AI-Frame/Brain/Defi-engine
ExecStart=/usr/bin/env python training/automation/bin/codex_event_watcher.py --queue training/automation/state/events.jsonl --rules training/automation/config/automation_rules.json --receipts-dir training/automation/receipts
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
